# ubuntu-automation/crontab-main.cron
#* * * * * /usr/bin/sudo /usr/bin/nr-update
*/10 * * * * /opt/newsreduce/crawler/src/main/bash/cleanup-empty-dirs.sh
*/5 * * * * /opt/newsreduce/crawler/src/main/bash/pull-remote-files.sh

